# roles/third_party/nginx-install/templates/nginx.service.j2
[Unit]
Description=The NGINX HTTP and reverse proxy server
After=syslog.target network-online.target remote-fs.target nss-lookup.target
Wants=network-online.target

[Service]
Type=forking
PIDFile={{ __nginx_pid_loc }}
ExecStartPre={{ __nginx_exec_path }} -t
ExecStart={{ __nginx_exec_path }}
ExecReload={{ __nginx_exec_path }} -s reload
ExecStop=/bin/kill -s QUIT $MAINPID
PrivateTmp=true

[Install]
WantedBy=multi-user.target