# roles/third_party/node-exporter-install/templates/etc/systemd/system/prometheus-node-exporter.service.j2
[Unit]
Description=Prometheus Node Exporter
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User={{ __prometheus_exporters_common_user }}
Group={{ __prometheus_exporters_common_group }}
ExecStart={{ __prometheus_exporters_common_root_dir }}/node_exporter_current/node_exporter --web.listen-address=:{{ __node_exporter_port }}

SyslogIdentifier=prometheus_node_exporter
Restart=always

[Install]
WantedBy=multi-user.target
